# meta-xt-prod-devel-rcar-control-gen4/recipes-guests/domu/files/domu.service
[Unit]
Description=DomU
Requires=backend-ready@bridge.service backend-ready@pcid.service
After=backend-ready@bridge.service backend-ready@pcid.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/xl create /etc/xen/domu.cfg
ExecStop=/usr/sbin/xl destroy DomU
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
